# default/udev/99-omarchy-platform-profile.rules
ACTION=="add", SUBSYSTEM=="platform-profile", ATTR{name}=="dell-pc", RUN+="/bin/chgrp wheel /sys/class/platform-profile/%k/profile"
ACTION=="add", SUBSYSTEM=="platform-profile", ATTR{name}=="dell-pc", RUN+="/bin/chmod g+w /sys/class/platform-profile/%k/profile"
