# config/systemd/user/omarchy-theme-auto.service
[Unit]
Description=Omarchy Auto Theme Switch
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/omarchy/bin/omarchy-theme-auto apply
Environment=DISPLAY=:0

