# config/systemd/user/omarchy-sleep-lock.service
[Unit]
Description=Lock Omarchy before suspend
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart=%h/.local/share/omarchy/bin/omarchy-system-sleep-monitor
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
