# config/systemd/user/omarchy-monitor-refresh-rate-watch.service
[Unit]
Description=Omarchy Monitor Refresh Rate Watcher
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/omarchy/bin/omarchy-monitor-refresh-rate-watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
