# config/systemd/user/omarchy-dynamic-bg.service
[Unit]
Description=Dynamic background rotation with awww
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/omarchy/bin/omarchy-dynamic-bg-rotate
