# config/systemd/user/omarchy-bg-cycle.timer
[Unit]
Description=Omarchy Background Cycle Timer

[Timer]
OnCalendar=hourly
AccuracySec=1min
Unit=omarchy-bg-cycle.service

[Install]
WantedBy=timers.target
