# config/systemd/user/omarchy-theme-schedule.timer
[Unit]
Description=Omarchy Theme Schedule Timer
Requires=omarchy-theme-schedule.service

[Timer]
# OnCalendar entries are written as a drop-in by omarchy-theme-schedule-apply
# at each firing. The drop-in lists today's and tomorrow's sunrise/sunset.
#
# Resume from suspend is handled by omarchy-theme-schedule-sleepwatch
# (a D-Bus listener on logind's PrepareForSleep signal), not by Persistent=true
# here — catching up missed events on resume raced with session stabilization
# and crashed Hyprland.
#
# Timezone changes are handled by omarchy-theme-schedule-tzwatch.path.
AccuracySec=30sec

[Install]
WantedBy=timers.target
