# config/systemd/user/omarchy-theme-schedule-tzwatch.path
[Unit]
Description=Watch /etc/localtime for timezone changes and reapply theme schedule

[Path]
PathChanged=/etc/localtime
Unit=omarchy-theme-schedule.service

[Install]
WantedBy=default.target
