# config/systemd/user/tzupdate.service
[Unit]
Description=Update system timezone from geoIP
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tzupdate
