# default/systemd/user/omarchy-update-user-notify.path
[Unit]
Description=Watch for Omarchy migrations

[Path]
PathExistsGlob=/usr/share/omarchy/migrations/*.sh
PathModified=/usr/share/omarchy/migrations
Unit=omarchy-update-user-notify.service

[Install]
WantedBy=graphical-session.target
