# default/systemd/user/omarchy-update-user-notify.service
[Unit]
Description=Notify about pending Omarchy user migrations
ConditionPathExistsGlob=/usr/share/omarchy/migrations/user/*.sh

[Service]
Type=oneshot
ExecStart=/usr/bin/omarchy-migrate-notify
