# default/systemd/user/omarchy-update-user-notify.service
[Unit]
Description=Notify about pending Omarchy migrations
ConditionPathIsDirectory=/usr/share/omarchy/migrations

[Service]
Type=oneshot
ExecStart=/usr/bin/omarchy-migrate-notify

[Install]
WantedBy=graphical-session.target
