# systemd/user/omarchy-theme-auto-switcher.timer
[Unit]
Description=Run Auto Theme Switcher periodically

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
