# systemd/user/omarchy-theme-auto-switcher.service
[Unit]
Description=Theme Auto Switcher Service

[Service]
Type=oneshot
ExecStart=%h/.local/share/omarchy/bin/omarchy-theme-next
